# docker/build/startup-script.conf
[program:startup-scripts]
process_name=%(program_name)s_%(process_num)02d
command=sh /scripts/startup.sh
autostart=true
autorestart=false
stopasgroup=true
killasgroup=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/error.log
stopwaitsecs=3600